# framework/files/etc/systemd/system/rancherd.service.d/override.conf
[Unit]
ConditionPathExists=!/run/cos/live_mode
ConditionPathExists=!/run/cos/recovery_mode
After=
After=cos-setup-network.service

[Service]
ExecStartPre=-/usr/bin/mkdir -p /usr/local/bin
