# framework/files/etc/systemd/system/elemental-install.service.d/override.conf
[Unit]
After=
After=cos-setup-network.service